# dune-weaver-touch/dune-weaver-touch.service
[Unit]
Description=Dune Weaver Touch Interface
After=multi-user.target
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/dune-weaver-touch

# EGLFS platform - direct framebuffer, no X11/Wayland needed
Environment=QT_QPA_PLATFORM=eglfs
Environment=QT_QPA_EGLFS_ALWAYS_SET_MODE=1
Environment=QT_QPA_EGLFS_PHYSICAL_WIDTH=155
Environment=QT_QPA_EGLFS_PHYSICAL_HEIGHT=86
Environment=QT_IM_MODULE=qtvirtualkeyboard

# Kiosk mode flag
Environment=KIOSK_MODE=1

# Wait for backend to be available before starting
ExecStartPre=/bin/bash -c 'until curl -s http://localhost:8080/serial_status > /dev/null 2>&1; do sleep 2; done'
ExecStart=/home/pi/dune-weaver-touch/venv/bin/python /home/pi/dune-weaver-touch/main.py

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target