# dune-weaver-touch/dune-weaver-touch.service
[Unit]
Description=Dune Weaver Touch Interface
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/dune-weaver-touch
Environment=DISPLAY=:0
Environment=QT_QPA_PLATFORM=linuxfb
Environment=QT_QPA_FB_DRM=1
Environment=QT_QPA_FONTDIR=/usr/share/fonts
ExecStart=/home/pi/dune-weaver-touch/venv/bin/python /home/pi/dune-weaver-touch/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target