# dune-weaver-touch/dune-weaver-touch.service
[Unit]
Description=Dune Weaver Touch Interface
After=multi-user.target graphical.target systemd-user-sessions.service
Wants=graphical.target
# Ensure framebuffer is ready before starting
After=systemd-vconsole-setup.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/dune-weaver-touch
Environment=DISPLAY=:0
Environment=QT_QPA_PLATFORM=linuxfb
Environment=QT_QPA_FB_DRM=1
Environment=QT_QPA_FONTDIR=/usr/share/fonts
# Use wrapper script that checks for framebuffer availability
ExecStart=/home/pi/dune-weaver-touch/start-with-fb-check.sh
Restart=on-failure
RestartSec=10
# Limit restart attempts on boot
StartLimitBurst=3
StartLimitIntervalSec=120
# Add detailed logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target