# scripts/dune-weaver-update.service
[Unit]
Description=Dune Weaver Update Watcher
After=multi-user.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/dune-weaver
ExecStart=/home/pi/dune-weaver/scripts/update-watcher.sh
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
